# idea-backlog/cron/com.matteisn.idea-backlog.weekly-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.matteisn.idea-backlog.weekly-digest</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/python3</string>
		<string>/Users/mattheweisner/.claude/skills/idea-backlog/scripts/weekly_digest.py</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/Users/mattheweisner/.claude/skills/idea-backlog/_workdir/cron-weekly-digest.log</string>
	<key>StandardOutPath</key>
	<string>/Users/mattheweisner/.claude/skills/idea-backlog/_workdir/cron-weekly-digest.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
		<key>Weekday</key>
		<integer>0</integer>
	</dict>
</dict>
</plist>
