# gtm-hub/launchd/com.matteisn.gtm-hub-regenerate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.matteisn.gtm-hub-regenerate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/mattheweisner/.claude/fakematt-today/gtm_hub_regenerate.py</string>
  </array>

  <!-- Every 15 minutes — index/decisions/README rebuild is idempotent + cheap -->
  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/mattheweisner/.claude/fakematt-today/gtm_hub.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mattheweisner/.claude/fakematt-today/gtm_hub.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/mattheweisner</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
