# gtm-hub/launchd/com.matteisn.gtm-hub-post.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.matteisn.gtm-hub-post</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/mattheweisner/.claude/fakematt-today/gtm_hub_post.py</string>
  </array>

  <!-- Monday 7:15 AM PT — weekly Slack digest -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/mattheweisner/.claude/fakematt-today/gtm_hub.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mattheweisner/.claude/fakematt-today/gtm_hub.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/mattheweisner</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
